# deploy/duckdns-caddy/notrus-duckdns.service
[Unit]
Description=Update DuckDNS for Notrus relay
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
EnvironmentFile=/etc/notrus/duckdns.env
ExecStart=/usr/local/bin/notrus-duckdns-update
